# pkgfiles/personal/libvirt-meta/.install
post_install() {
  echo -e "\nTo enable the libvirt service execute the following commands with root permission:\n"

  echo "systemctl enable --now libvirtd"
  echo "virsh net-start default"
  echo "systemctl restart libvirtd"
  echo "usermod -G libvirt -a $USER"
}
